# tools/gh-agent-bridge/launchd/com.heritagedx.ghbridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.heritagedx.ghbridge</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/heesungbae/.nvm/versions/node/v22.22.0/bin/node</string>
    <string>/Users/heesungbae/heritage-dx/tools/gh-agent-bridge/node_modules/.bin/tsx</string>
    <string>/Users/heesungbae/heritage-dx/tools/gh-agent-bridge/src/index.ts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/heesungbae/heritage-dx/tools/gh-agent-bridge</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/heesungbae/.nvm/versions/node/v22.22.0/bin</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/heesungbae/Library/Logs/gh-agent-bridge.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/heesungbae/Library/Logs/gh-agent-bridge.err.log</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
